# crontab/cronfile
# see http://www.adminschoice.com/docs/crontab.htm
#
# to install a cron job, type 'crontab -u username cronfile'
# once the cron job is installed, edit it with 'crontab -u username -e'
#
# m h  dom mon dow   command
0 4 * * * /usr/bin/python /var/virtualhost/sites/feedme/prod/server/term_vector.py reindex | mail feedme@csail.mit.edu -s "Nightly Term Vector Script"
* 0-3,5-23 * * * /usr/bin/python /var/virtualhost/sites/feedme/prod/server/term_vector.py incremental >/dev/null 2>&1
0 15 * * 2,5 /usr/bin/python /var/virtualhost/sites/feedme/prod/server/digest.py | mail feedme@csail.mit.edu -s "Digest script output"
